# scripts/com.globalapex.governance-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  RETIRED 2026-04-18. DO NOT LOAD.
  OpenAI plans cancelled, API keys revoked, all 10 Paperclip agents paused.
  Embedded Postgres on :54329 no longer running.
  Kept as historical reference. See docs/09-Paperclip-Retirement-Handover-2026-04-18.md
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.globalapex.governance-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/mmh/.paperclip/governance-sync.py</string>
        <string>--report</string>
    </array>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/governance-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/governance-sync.log</string>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
